# deploy/systemd/scalpx-mme-close.service
[Unit]
Description=Stop ScalpX MME Canonical Runtime

[Service]
Type=oneshot
User=root
ExecStart=/bin/systemctl stop scalpx-mme.service
