# deploy/systemd/scalpx-mme.service
[Unit]
Description=ScalpX MME Canonical Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=Lenovo
WorkingDirectory=/home/Lenovo/scalpx/projects/mme_scalpx
Environment=MME_BOOTSTRAP_PROVIDER=app.mme_scalpx.integrations.bootstrap_provider:provide
ExecStart=/home/Lenovo/scalpx/projects/mme_scalpx/.venv/bin/python -m app.mme_scalpx.main
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
