# deploy/systemd/scalpx-mme-open.timer
[Unit]
Description=Start ScalpX MME at 08:55 IST on weekdays

[Timer]
OnCalendar=Mon..Fri 08:55:00 Asia/Kolkata
Persistent=true
Unit=scalpx-mme-open.service

[Install]
WantedBy=timers.target
