# deploy/systemd/scalpx-mme-close.timer
[Unit]
Description=Stop ScalpX MME at 15:35 IST on weekdays

[Timer]
OnCalendar=Mon..Fri 15:35:00 Asia/Kolkata
Persistent=true
Unit=scalpx-mme-close.service

[Install]
WantedBy=timers.target
