# deploy/systemd/scalpx-mme-open.service
[Unit]
Description=Start ScalpX MME Canonical Runtime

[Service]
Type=oneshot
User=root
ExecStart=/bin/systemctl start scalpx-mme.service
